# platforms/epel/cryptoded.service
[Unit]
Description=Cryptode-Hidden Path
After=syslog.target network.target

[Service]
ExecStart=/usr/sbin/cryptoded -D

ExecStop=kill -TERM `cat /var/run/cryptoded/cryptoded.pid`
ExecReload=kill -SIGUSR1 `cat /var/run/cryptoded/cryptoded.pid`

Restart=always

Type=notify
NotifyAccess=main

[Install]
WantedBy=multi-user.target
